# deploy/whatsapp-worker.service
[Unit]
Description=WhatsApp Worker (pure-API)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=wa-worker
WorkingDirectory=/home/wa-worker/app
ExecStart=/usr/bin/node --enable-source-maps --import tsx src/index.ts
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/wa-worker/worker.log
StandardError=append:/var/log/wa-worker/worker.log

# --- Security hardening ---
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/home/wa-worker/app /var/lib/wa-worker /var/log/wa-worker
PrivateTmp=yes
ProtectHome=read-only
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=yes
# Chromium needs writable+executable pages for V8's JIT.
MemoryDenyWriteExecute=no
RestrictRealtime=yes
SystemCallArchitectures=native

# --- Resource limits (Chromium can be greedy) ---
MemoryHigh=2G
MemoryMax=3G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
